# host/systemd/meshcore-hostbot.service
[Service]
ExecStart=/usr/bin/python3 /opt/meshcore-hostbot/host_metrics.py
