# host/systemd/meshcore-host-action.path
[Path]
PathExists=/opt/meshcore-hostbot/requests/host_action.json
