# host/systemd/meshcore-hostbot.service
[Unit]
Description=Write host metrics for MeshCore RemoteTerm bot

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/meshcore-hostbot/host_metrics.py
