# host/systemd/meshcore-host-action.service
[Service]
ExecStart=/usr/bin/python3 /opt/meshcore-hostbot/handle_host_request.py
